# net-firewall/bango/files/init
#!/sbin/openrc-run
#extra_commands="publish"

start() {
	ebegin "Starting ${SVCNAME}"
	rm -f /var/run/bango.pid
	start-stop-daemon --start -b --exec "/usr/bin/bango" \
	    -m --pidfile "/var/run/bango.pid" -1  /var/log/bango.log -2  /var/log/bango.log
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --exec "/usr/bin/bango" \
	    --pidfile "/var/run/bango.pid" --quiet
	eend $?
}

