# dev-db/dragonfly/files/dragonfly.initd
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/dragonfly"
command_args="${DRAGONFLY_OPTS}"
command_background="true"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	use localmount logger
	after keepalived
}
