# dev-java/archiva-bin/files/init.d/archiva-bin-1.1.3
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	use dns logger net
}

start()	{
	ebegin "Starting Apache Archiva"

	# why?
	export ARCHIVA_HOME=$ARCHIVA_HOME
	export ARCHIVA_BASE=$ARCHIVA_BASE

	JAVA_HOME=`java-config --jre-home`
	CLASSPATH="${CLASSPATH}:${ARCHIVA_HOME}/lib:${JAVA_HOME}/lib/tools.jar"

#		--chuid ${ARCHIVA_USER}:${ARCHIVA_GROUP} \

	start-stop-daemon --start --quiet \
		--make-pidfile --pidfile /var/run/archiva-bin.pid --exec \
		"${ARCHIVA_BIN}" ${ARCHIVA_ARGS} >> "${ARCHIVA_BASE}"/logs/archiva.out 2>&1 &
	eend $?
}

stop()	{
	ebegin "Stopping Apache Archiva"
	start-stop-daemon --stop --quiet --pidfile /var/run/archiva-bin.pid
	eend $?
}
