# sys-cluster/saunafs/files/saunafscgiserver.initd
#!/sbin/openrc-run

PIDFILE=/var/run/sfscgiserv.pid

depend() {
	need net
}

start() {
	ebegin "Starting sfs CGI server"
	start-stop-daemon --start --quiet --exec /usr/sbin/sfscgiserv \
	--pidfile "${PIDFILE}" --make-pidfile --background \
	-- -H ${BIND_HOST} -P ${BIND_PORT} -f
	eend $?
}

stop() {
	ebegin "Stopping sfs CGI server"
	start-stop-daemon --stop -q --pidfile "${PIDFILE}"
	eend $?
}
